# ops/systemd/qnty-data-refresh.service
[Unit]
Description=Qnty Data Refresh - Fetch OHLCV and Funding from Binance Public API
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=qnty
Group=qnty
WorkingDirectory=/srv/qnty/repo
Environment="PATH=/srv/qnty/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="VIRTUAL_ENV=/srv/qnty/venv"
ExecStart=/srv/qnty/repo/ops/bin/qnty-data-refresh.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=qnty-data-refresh
# Timeout after 10 minutes
TimeoutSec=600
