# ops/systemd/qnty-shadow-run.service
[Unit]
Description=Qnty Shadow Run - Frozen Package V2 Observer on New 8h Bar
After=network.target
Wants=network.target

[Service]
Type=oneshot
User=qnty
Group=qnty
WorkingDirectory=/srv/qnty/repo
Environment="PATH=/srv/qnty/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="VIRTUAL_ENV=/srv/qnty/venv"
ExecStart=/srv/qnty/repo/ops/bin/qnty-shadow-run.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=qnty-shadow-run
# Timeout after 30 minutes (stage4 + validation can take time)
TimeoutSec=1800
